# etc/s6-init/run-image/service/early-getty-ttyS0/run
#!/bin/execlineb -P
s6-envdir env
multisubstitute {
	import -i HOSTNAME
	import -iu TTY
	import -iu BAUDS
}
cd /
/sbin/getty -H $HOSTNAME -n -L /bin/ash $BAUDS $TTY
