# etc/s6-init/run-image/service/s6-linux-init-early-getty/run
#!/bin/execlineb -P
cd /
/sbin/getty -H alpy -n -l "/bin/ash" ttyS0 115200
